# deploy/hannah-webui.service
[Unit]
Description=Hannah WebUI
Documentation=https://dev.kernstock.net/gessinger/voice/hannah
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hannah
WorkingDirectory=/opt/hannah/webui
Environment=HANNAH_WEBUI_CONFIG=/etc/hannah-webui/config.yaml
# HOME zeigt sonst auf /home/hannah, das ProtectHome=true unsichtbar macht — grpc's
# C-Core scheitert beim Fork-Handling der gunicorn-Worker daran mit "Permission
# denied" (#104). Gleiches Muster wie bei voiceid (siehe dessen .service).
Environment=HOME=/opt/hannah/webui
# Bind-Adresse ist hier festgelegt, nicht in config.yaml — gunicorn bindet den
# Socket selbst, bevor die WSGI-App (wsgi.py) überhaupt geladen wird. config.yaml's
# host/port gelten nur für main.py's Flask-Dev-Server (lokale Entwicklung).
ExecStart=/opt/hannah/webui/venv/bin/gunicorn --workers 2 --bind 0.0.0.0:5000 wsgi:app
Restart=on-failure
RestartSec=5s

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/etc/hannah-webui
ReadWritePaths=/opt/hannah/webui

[Install]
WantedBy=multi-user.target
